# src/service/hello.service
[Unit]
Description=HelloWorld

[Service]
ExecStart=/usr/local/hello_world/hello_world.sh

[Install]
WantedBy=multi-user.target
